# deploy/tutor.service.example
# ============================================================================
# systemd unit для tutor_site (gunicorn).
# Положить в /etc/systemd/system/tutor.service, затем:
#   systemctl daemon-reload
#   systemctl enable --now tutor
#   systemctl status tutor
#   journalctl -u tutor -f
# ============================================================================

[Unit]
Description=Tutor site (gunicorn)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=notify
User=tutor
Group=tutor
WorkingDirectory=/opt/tutor
EnvironmentFile=/opt/tutor/.env

# unix-сокет для общения с nginx (быстрее tcp).
# Сокет лежит ВНУТРИ /run/tutor/, потому что ProtectSystem=strict делает
# сам /run/ read-only — а /run/tutor/ создаётся через RuntimeDirectory ниже.
ExecStart=/opt/tutor/venv/bin/gunicorn \
    --workers 3 \
    --worker-class sync \
    --timeout 60 \
    --bind unix:/run/tutor/tutor.sock \
    --umask 007 \
    --access-logfile - \
    --error-logfile - \
    tutor_core.wsgi:application

# Перезапускать при крашах.
Restart=on-failure
RestartSec=3

# systemd создаёт /run/tutor/ с правами `tutor:tutor`.
# Чтобы nginx (от www-data) мог стучаться в сокет, добавить www-data в группу tutor:
#   usermod -aG tutor www-data && systemctl reload nginx
RuntimeDirectory=tutor
RuntimeDirectoryMode=0755

# Лимиты безопасности.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
# Не включаем /opt/tutor/db.sqlite3 — если у тебя PostgreSQL, файла нет
# и systemd упадёт с status=226/NAMESPACE.
ReadWritePaths=/opt/tutor/media /opt/tutor/staticfiles /run/tutor

[Install]
WantedBy=multi-user.target
